# src/container-manager/docker-files/rif-comms-bootnode/assets/supervisord.conf
[supervisord]
nodaemon=true

[program:rif-comms]
command=/usr/local/bin/startRifComms
autostart=true
user=root
startsecs=0
startretries=0
exitcodes=0
